# services/timelapse-viewer.service
[Unit]
Description=Aquarium Timelapse Web Viewer
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/saainithil97/projects/timelapse
ExecStart=/usr/bin/python3 /home/saainithil97/projects/timelapse/viewer.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
